# roles/apps/node-exporter/templates/etc/systemd/system/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User={{ node_exporter_user }}
Group={{ node_exporter_user }}
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=multi-user.target