# roles/apps/redis-exporter/templates/etc/systemd/system/redis_exporter.service.j2
[Unit]
Description=Redis Exporter

[Service]
User={{ redis_exporter_user }}
Group={{ redis_exporter_user }}
ExecStart=/usr/local/bin/redis_exporter

[Install]
WantedBy=multi-user.target